# crontab
0 * * 2 *    sh ~/MalodyChartJam/cron.sh active
0 0 * 3 *    sh ~/MalodyChartJam/cron.sh after
